# scripts/cron-meshctl/etc/crontabs/root
*/1 * * * * /bin/bash "/root/cron-meshctl.sh" >/dev/null 2>&1
